# deploy/platform/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# =============================================================================
# Infrastructure services (start first — no Kafka dependency)
# =============================================================================

[program:keycloak]
command=/opt/keycloak/bin/kc.sh start-dev
    --http-port=%(ENV_KEYCLOAK_PORT)s
    --import-realm
environment=
    KEYCLOAK_ADMIN=%(ENV_KEYCLOAK_ADMIN)s,
    KEYCLOAK_ADMIN_PASSWORD=%(ENV_KEYCLOAK_ADMIN_PASSWORD)s,
    KC_DB=dev-file,
    KC_HOSTNAME_STRICT=false,
    KC_HOSTNAME_STRICT_HTTPS=false,
    KC_HTTP_ENABLED=true
autostart=true
autorestart=true
startretries=3
startsecs=15
priority=10
stdout_logfile=/var/log/supervisor/keycloak.log
stderr_logfile=/var/log/supervisor/keycloak.err
stdout_logfile_maxbytes=20MB
redirect_stderr=false

# =============================================================================
# Ingestion services (require Kafka — started after wait-for-kafka)
# =============================================================================

[program:dcim-ingestion]
command=/app/dcim-ingestion/dcim-ingestion-service
    --port=%(ENV_DCIM_INGESTION_PORT)s
    --kafka=%(ENV_KAFKA_BOOTSTRAP)s
    --pg-host=%(ENV_POSTGRES_HOST)s
    --pg-port=%(ENV_POSTGRES_PORT)s
    --pg-db=%(ENV_POSTGRES_DB)s
    --pg-user=%(ENV_POSTGRES_USER)s
    --pg-pass=%(ENV_POSTGRES_PASSWORD)s
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=50
stdout_logfile=/var/log/supervisor/dcim-ingestion.log
stderr_logfile=/var/log/supervisor/dcim-ingestion.err
stdout_logfile_maxbytes=20MB

[program:telemetry-gateway]
command=/app/telemetry-gateway/telemetry-gateway-service
    --port=%(ENV_TELEMETRY_GATEWAY_PORT)s
    --udp-port=%(ENV_TELEMETRY_GATEWAY_UDP_PORT)s
    --kafka=%(ENV_KAFKA_BOOTSTRAP)s
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=50
stdout_logfile=/var/log/supervisor/telemetry-gateway.log
stderr_logfile=/var/log/supervisor/telemetry-gateway.err
stdout_logfile_maxbytes=20MB

[program:active-discovery]
command=/app/active-discovery/active-discovery-service
    --port=%(ENV_ACTIVE_DISCOVERY_PORT)s
    --kafka=%(ENV_KAFKA_BOOTSTRAP)s
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=50
stdout_logfile=/var/log/supervisor/active-discovery.log
stderr_logfile=/var/log/supervisor/active-discovery.err
stdout_logfile_maxbytes=20MB

# =============================================================================
# Processing services (require Kafka + stores)
# =============================================================================

[program:graph-updater]
command=/app/graph-updater/graph-updater-service
    --port=%(ENV_GRAPH_UPDATER_PORT)s
    --kafka=%(ENV_KAFKA_BOOTSTRAP)s
    --neo4j=bolt://%(ENV_NEO4J_HOST)s:%(ENV_NEO4J_BOLT_PORT)s
    --neo4j-user=%(ENV_NEO4J_USER)s
    --neo4j-pass=%(ENV_NEO4J_PASSWORD)s
    --pg-host=%(ENV_POSTGRES_HOST)s
    --pg-db=%(ENV_POSTGRES_DB)s
    --pg-user=%(ENV_POSTGRES_USER)s
    --pg-pass=%(ENV_POSTGRES_PASSWORD)s
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=60
stdout_logfile=/var/log/supervisor/graph-updater.log
stderr_logfile=/var/log/supervisor/graph-updater.err
stdout_logfile_maxbytes=20MB

[program:timeseries-writer]
command=/app/timeseries-writer/timeseries-writer-service
    --port=%(ENV_TIMESERIES_WRITER_PORT)s
    --kafka=%(ENV_KAFKA_BOOTSTRAP)s
    --ts-host=%(ENV_TIMESCALE_HOST)s
    --ts-port=%(ENV_TIMESCALE_PORT)s
    --ts-db=%(ENV_TIMESCALE_DB)s
    --ts-user=%(ENV_TIMESCALE_USER)s
    --ts-pass=%(ENV_TIMESCALE_PASSWORD)s
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=60
stdout_logfile=/var/log/supervisor/timeseries-writer.log
stderr_logfile=/var/log/supervisor/timeseries-writer.err
stdout_logfile_maxbytes=20MB

[program:event-archive]
command=/app/event-archive/event-archive-service
    --port=%(ENV_EVENT_ARCHIVE_PORT)s
    --kafka=%(ENV_KAFKA_BOOTSTRAP)s
    --s3-endpoint=%(ENV_S3_ENDPOINT)s
    --s3-key=%(ENV_S3_ACCESS_KEY)s
    --s3-secret=%(ENV_S3_SECRET_KEY)s
    --s3-bucket=%(ENV_S3_BUCKET)s
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=60
stdout_logfile=/var/log/supervisor/event-archive.log
stderr_logfile=/var/log/supervisor/event-archive.err
stdout_logfile_maxbytes=20MB

# =============================================================================
# Health aggregator — single /health endpoint for the whole container
# =============================================================================

[program:health-aggregator]
command=python3 /app/health-aggregator/main.py
    --port=8099
    --services=%(ENV_DCIM_INGESTION_PORT)s,%(ENV_TELEMETRY_GATEWAY_PORT)s,%(ENV_ACTIVE_DISCOVERY_PORT)s,%(ENV_GRAPH_UPDATER_PORT)s,%(ENV_TIMESERIES_WRITER_PORT)s,%(ENV_EVENT_ARCHIVE_PORT)s
autostart=true
autorestart=true
priority=90
stdout_logfile=/var/log/supervisor/health-aggregator.log
stderr_logfile=/var/log/supervisor/health-aggregator.err
