# deploy/agents/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
loglevel=info
pidfile=/var/run/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# =============================================================================
# MLflow tracking server — starts first (agents register models to it)
# =============================================================================

[program:mlflow]
command=mlflow server
    --backend-store-uri sqlite:////mlflow/db/mlflow.db
    --default-artifact-root /mlflow/artifacts
    --host 0.0.0.0
    --port %(ENV_MLFLOW_PORT)s
autostart=true
autorestart=true
startretries=3
startsecs=5
priority=10
stdout_logfile=/var/log/supervisor/mlflow.log
stderr_logfile=/var/log/supervisor/mlflow.err
stdout_logfile_maxbytes=20MB

# =============================================================================
# topology-reconciliation-agent (Python FastAPI)
# Rule-based graph comparison — no heavy ML in Phase 1.
# Publishes to drift.events and drift.suggestions.
# =============================================================================

[program:topology-agent]
command=uvicorn main:app
    --host 0.0.0.0
    --port %(ENV_TOPOLOGY_AGENT_PORT)s
    --workers 2
    --log-level warning
directory=/app/topology-agent
environment=
    NEO4J_URI=bolt://%(ENV_NEO4J_HOST)s:%(ENV_NEO4J_BOLT_PORT)s,
    NEO4J_USER=%(ENV_NEO4J_USER)s,
    NEO4J_PASSWORD=%(ENV_NEO4J_PASSWORD)s,
    KAFKA_BOOTSTRAP=%(ENV_KAFKA_BOOTSTRAP)s,
    POSTGRES_HOST=%(ENV_POSTGRES_HOST)s,
    POSTGRES_PORT=%(ENV_POSTGRES_PORT)s,
    POSTGRES_DB=%(ENV_POSTGRES_DB)s,
    POSTGRES_USER=%(ENV_POSTGRES_USER)s,
    POSTGRES_PASSWORD=%(ENV_POSTGRES_PASSWORD)s,
    MLFLOW_TRACKING_URI=http://localhost:%(ENV_MLFLOW_PORT)s
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=50
stdout_logfile=/var/log/supervisor/topology-agent.log
stderr_logfile=/var/log/supervisor/topology-agent.err
stdout_logfile_maxbytes=20MB

# =============================================================================
# classification-capability-agent (Python FastAPI + LightGBM)
# Device type inference from SNMP/Redfish fingerprints.
# Publishes to classification.results.
# =============================================================================

[program:classification-agent]
command=uvicorn main:app
    --host 0.0.0.0
    --port %(ENV_CLASSIFICATION_AGENT_PORT)s
    --workers 2
    --log-level warning
directory=/app/classification-agent
environment=
    NEO4J_URI=bolt://%(ENV_NEO4J_HOST)s:%(ENV_NEO4J_BOLT_PORT)s,
    NEO4J_USER=%(ENV_NEO4J_USER)s,
    NEO4J_PASSWORD=%(ENV_NEO4J_PASSWORD)s,
    KAFKA_BOOTSTRAP=%(ENV_KAFKA_BOOTSTRAP)s,
    POSTGRES_HOST=%(ENV_POSTGRES_HOST)s,
    POSTGRES_PORT=%(ENV_POSTGRES_PORT)s,
    POSTGRES_DB=%(ENV_POSTGRES_DB)s,
    POSTGRES_USER=%(ENV_POSTGRES_USER)s,
    POSTGRES_PASSWORD=%(ENV_POSTGRES_PASSWORD)s,
    MLFLOW_TRACKING_URI=http://localhost:%(ENV_MLFLOW_PORT)s
autostart=true
autorestart=true
startretries=5
startsecs=15
priority=60
stdout_logfile=/var/log/supervisor/classification-agent.log
stderr_logfile=/var/log/supervisor/classification-agent.err
stdout_logfile_maxbytes=20MB
