# deploy/api/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
loglevel=info
pidfile=/var/run/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# =============================================================================
# nginx — starts first, upstream services may be temporarily unavailable
# nginx handles this gracefully (502 until upstreams are ready)
# =============================================================================

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
startsecs=2
priority=10
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
stdout_logfile_maxbytes=20MB

# =============================================================================
# API services — NestJS (Node.js)
# =============================================================================

[program:graph-api]
command=node /app/graph-api/dist/main.js
environment=
    PORT=%(ENV_GRAPH_API_INTERNAL_PORT)s,
    NEO4J_URI=bolt://%(ENV_NEO4J_HOST)s:%(ENV_NEO4J_BOLT_PORT)s,
    NEO4J_USER=%(ENV_NEO4J_USER)s,
    NEO4J_PASSWORD=%(ENV_NEO4J_PASSWORD)s,
    TIMESCALE_HOST=%(ENV_TIMESCALE_HOST)s,
    TIMESCALE_PORT=%(ENV_TIMESCALE_PORT)s,
    TIMESCALE_DB=%(ENV_TIMESCALE_DB)s,
    TIMESCALE_USER=%(ENV_TIMESCALE_USER)s,
    TIMESCALE_PASSWORD=%(ENV_TIMESCALE_PASSWORD)s,
    KEYCLOAK_URL=http://flowcore-platform:%(ENV_KEYCLOAK_PORT)s,
    KEYCLOAK_REALM=%(ENV_KEYCLOAK_REALM)s,
    JWT_SECRET=%(ENV_JWT_SECRET)s
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=50
stdout_logfile=/var/log/supervisor/graph-api.log
stderr_logfile=/var/log/supervisor/graph-api.err
stdout_logfile_maxbytes=20MB

[program:insights-api]
command=node /app/insights-api/dist/main.js
environment=
    PORT=%(ENV_INSIGHTS_API_INTERNAL_PORT)s,
    POSTGRES_HOST=%(ENV_POSTGRES_HOST)s,
    POSTGRES_PORT=%(ENV_POSTGRES_PORT)s,
    POSTGRES_DB=%(ENV_POSTGRES_DB)s,
    POSTGRES_USER=%(ENV_POSTGRES_USER)s,
    POSTGRES_PASSWORD=%(ENV_POSTGRES_PASSWORD)s,
    NEO4J_URI=bolt://%(ENV_NEO4J_HOST)s:%(ENV_NEO4J_BOLT_PORT)s,
    NEO4J_USER=%(ENV_NEO4J_USER)s,
    NEO4J_PASSWORD=%(ENV_NEO4J_PASSWORD)s,
    KAFKA_BOOTSTRAP=%(ENV_KAFKA_BOOTSTRAP)s,
    KEYCLOAK_URL=http://flowcore-platform:%(ENV_KEYCLOAK_PORT)s,
    KEYCLOAK_REALM=%(ENV_KEYCLOAK_REALM)s,
    JWT_SECRET=%(ENV_JWT_SECRET)s
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=50
stdout_logfile=/var/log/supervisor/insights-api.log
stderr_logfile=/var/log/supervisor/insights-api.err
stdout_logfile_maxbytes=20MB

[program:eventing-integration]
command=node /app/eventing-integration/dist/main.js
environment=
    PORT=%(ENV_EVENTING_INTEGRATION_INTERNAL_PORT)s,
    KAFKA_BOOTSTRAP=%(ENV_KAFKA_BOOTSTRAP)s,
    POSTGRES_HOST=%(ENV_POSTGRES_HOST)s,
    POSTGRES_PORT=%(ENV_POSTGRES_PORT)s,
    POSTGRES_DB=%(ENV_POSTGRES_DB)s,
    POSTGRES_USER=%(ENV_POSTGRES_USER)s,
    POSTGRES_PASSWORD=%(ENV_POSTGRES_PASSWORD)s,
    JWT_SECRET=%(ENV_JWT_SECRET)s
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=50
stdout_logfile=/var/log/supervisor/eventing-integration.log
stderr_logfile=/var/log/supervisor/eventing-integration.err
stdout_logfile_maxbytes=20MB

# =============================================================================
# NOC frontend — Vite dev server (hot reload in development)
# In production: replace with `nginx serve /app/noc-frontend/dist`
# =============================================================================

[program:noc-frontend]
command=npm run dev -- --host 0.0.0.0 --port %(ENV_NOC_FRONTEND_INTERNAL_PORT)s
directory=/app/noc-frontend
environment=
    VITE_GRAPHQL_URL=http://localhost:%(ENV_API_GATEWAY_PORT)s/graphql,
    VITE_INSIGHTS_URL=http://localhost:%(ENV_API_GATEWAY_PORT)s/api/insights,
    VITE_KEYCLOAK_URL=http://localhost:%(ENV_KEYCLOAK_PORT)s,
    VITE_KEYCLOAK_REALM=%(ENV_KEYCLOAK_REALM)s,
    VITE_KEYCLOAK_CLIENT_ID=flowcore-ui
autostart=true
autorestart=true
startretries=3
startsecs=15
priority=60
stdout_logfile=/var/log/supervisor/noc-frontend.log
stderr_logfile=/var/log/supervisor/noc-frontend.err
stdout_logfile_maxbytes=20MB
