# docker/supervisord/agents.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=50MB
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:mlflow]
command=mlflow server
  --backend-store-uri sqlite:////mlflow/db/mlflow.db
  --default-artifact-root /mlflow/artifacts
  --host 0.0.0.0
  --port %(ENV_MLFLOW_PORT)s
  --serve-artifacts
priority=10
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/supervisord/mlflow.log
stderr_logfile=/var/log/supervisord/mlflow.log

[program:topology-agent]
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_TOPO_AGENT_PORT)s --workers 1 --log-level %(ENV_LOG_LEVEL)s
directory=/app/topology-reconciliation-agent
priority=20
autostart=true
autorestart=true
startretries=10
startsecs=15
stdout_logfile=/var/log/supervisord/topology-agent.log
stderr_logfile=/var/log/supervisord/topology-agent.log

[program:classification-agent]
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_CLASS_AGENT_PORT)s --workers 1 --log-level %(ENV_LOG_LEVEL)s
directory=/app/classification-capability-agent
priority=21
autostart=true
autorestart=true
startretries=10
startsecs=15
stdout_logfile=/var/log/supervisord/classification-agent.log
stderr_logfile=/var/log/supervisord/classification-agent.log
