# docker/supervisord/api-ui.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=50MB
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command=nginx -g "daemon off;"
priority=5
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/var/log/supervisord/nginx.log
stderr_logfile=/var/log/supervisord/nginx.log

[program:graph-api]
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_GRAPH_API_INTERNAL_PORT)s --log-level %(ENV_LOG_LEVEL)s
directory=/app/graph-api
priority=20
autostart=true
autorestart=true
startretries=10
startsecs=10
stdout_logfile=/var/log/supervisord/graph-api.log
stderr_logfile=/var/log/supervisord/graph-api.log

[program:insights-api]
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_INSIGHTS_API_INTERNAL_PORT)s --workers 2 --log-level %(ENV_LOG_LEVEL)s
directory=/app/insights-api
priority=21
autostart=true
autorestart=true
startretries=10
startsecs=10
stdout_logfile=/var/log/supervisord/insights-api.log
stderr_logfile=/var/log/supervisord/insights-api.log

[program:eventing-integration]
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_EVENTING_API_INTERNAL_PORT)s --log-level %(ENV_LOG_LEVEL)s
directory=/app/eventing-integration
priority=22
autostart=true
autorestart=true
startretries=10
startsecs=10
stdout_logfile=/var/log/supervisord/eventing-integration.log
stderr_logfile=/var/log/supervisord/eventing-integration.log
