# docker/supervisord/platform.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=50MB
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:keycloak]
command=/opt/keycloak/bin/kc.sh start-dev --http-port=8080 --import-realm
environment=KEYCLOAK_ADMIN=%(ENV_KEYCLOAK_ADMIN)s,KEYCLOAK_ADMIN_PASSWORD=%(ENV_KEYCLOAK_ADMIN_PASSWORD)s,KC_DB=dev-mem,KC_HOSTNAME_STRICT=false,KC_HTTP_ENABLED=true
priority=10
autostart=true
autorestart=true
startretries=3
startsecs=20
stdout_logfile=/var/log/supervisord/keycloak.log
stderr_logfile=/var/log/supervisord/keycloak.log

[program:dcim-ingestion]
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_DCIM_INGESTION_PORT)s --log-level %(ENV_LOG_LEVEL)s
directory=/app/dcim-ingestion
priority=30
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/var/log/supervisord/dcim-ingestion.log
stderr_logfile=/var/log/supervisord/dcim-ingestion.log

[program:telemetry-gateway]
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_TELEMETRY_GATEWAY_HTTP_PORT)s --workers 2 --log-level %(ENV_LOG_LEVEL)s
directory=/app/telemetry-gateway
priority=31
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/var/log/supervisord/telemetry-gateway.log
stderr_logfile=/var/log/supervisord/telemetry-gateway.log

[program:active-discovery]
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_ACTIVE_DISCOVERY_PORT)s --log-level %(ENV_LOG_LEVEL)s
directory=/app/active-discovery
priority=32
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/var/log/supervisord/active-discovery.log
stderr_logfile=/var/log/supervisord/active-discovery.log

[program:graph-updater]
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_GRAPH_UPDATER_PORT)s --log-level %(ENV_LOG_LEVEL)s
directory=/app/graph-updater
priority=50
autostart=true
autorestart=true
startretries=10
startsecs=10
stdout_logfile=/var/log/supervisord/graph-updater.log
stderr_logfile=/var/log/supervisord/graph-updater.log

[program:timeseries-writer]
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_TIMESERIES_WRITER_PORT)s --log-level %(ENV_LOG_LEVEL)s
directory=/app/timeseries-writer
priority=51
autostart=true
autorestart=true
startretries=10
startsecs=10
stdout_logfile=/var/log/supervisord/timeseries-writer.log
stderr_logfile=/var/log/supervisord/timeseries-writer.log

[program:event-archive]
command=uvicorn main:app --host 0.0.0.0 --port %(ENV_EVENT_ARCHIVE_PORT)s --log-level %(ENV_LOG_LEVEL)s
directory=/app/event-archive
priority=52
autostart=true
autorestart=true
startretries=10
startsecs=10
stdout_logfile=/var/log/supervisord/event-archive.log
stderr_logfile=/var/log/supervisord/event-archive.log
